# audio/radiocapture-cleanup.cron
0 6 * * * root /opt/radiocapture-motorola/audio/cleanup.sh
